# dist/systemd/alkalinethundermud.service
[Unit]
Description=Alkaline Thunder MUD
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=np
ExecStart=/opt/socially-distant/server/AlkalineThunderMud.Server

[Install]
WantedBy=multi-user.target
